# web/supervisord.conf
[program:app]
command = python3 /iomt-project/mqtt-daemon/run.py

[program:mosquitto]
command = /usr/sbin/mosquitto -c /etc/mosquitto.conf
autorestart=true
